# sr/recipes/sr-udev-rules/files/99-firmware.rules
SUBSYSTEM=="firmware", ACTION=="add", RUN+="/bin/fw-prodder $env{FIRMWARE} $env{DEVPATH}"
