# sr/recipes/roboscripts/files/zzzrobot.rules
SUBSYSTEM=="block", ACTION=="add"    RUN+="/etc/robot/insert.sh"
SUBSYSTEM=="block", ACTION=="remove" RUN+="/etc/robot/remove.sh"
